# .config/systemd/user/duplicacy-backup.timer
[Unit]
Description=Run the duplicacy backup daily.

[Timer]
OnCalendar=daily
AccuracySec=10min
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
