# .config/systemd/user/low-battery.service
[Unit]
Description=Notify or suspend on low battery

[Service]
ExecStart=%h/.local/bin/low-battery-action
Type=oneshot
