# src/config/systemd/user/auto-suspend.service
[Unit]
Description=Suspend the OS

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl suspend
